# deploy/polymarket-arb-daemon.service
[Unit]
Description=Polymarket always-on websocket monotonic-arb daemon (Japan)
Documentation=https://github.com/virtualwriter/polymarket-trader
After=network-online.target
Wants=network-online.target
# Repeated crashes within 5 min -> stop trying so we don't burn the API / hammer
# the CLOB. Reset by a successful long-lived run.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/polymarket-trader
EnvironmentFile=-/etc/polymarket-pm-executor.env
EnvironmentFile=-/etc/polymarket-trader.env
Environment=POLYMARKET_TRADER_STATE_DIR=/var/lib/polymarket-trader
Environment=NODE_NO_WARNINGS=1
ExecStart=/usr/local/bin/run-polymarket-arb-daemon
Restart=always
RestartSec=10
# The daemon logs a heartbeat; if it wedges, systemd restarts it. The daemon
# does not implement sd_notify, so WatchdogSec is intentionally omitted to avoid
# false-positive kills. Liveness is covered by Restart=always + reconnect loops.
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
