# deploy/sports-arb-llm-learn.service
[Unit]
Description=Sports Arb LLM daily learning pass
After=network-online.target sports-arb-daemon.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/sports-arb
EnvironmentFile=/etc/sports-arb.env
Environment=SPORTS_ARB_STATE_DIR=/var/lib/sports-arb
Environment=SPORTS_ARB_DATA_DIR=/var/lib/sports-arb/data
Environment=SPORTS_ARB_RUNTIME_DIR=/var/lib/sports-arb/runtime
Environment=NODE_NO_WARNINGS=1
ExecStart=/bin/bash /opt/sports-arb/scripts/nightly-digest.sh llm
User=root
Group=root
TimeoutStartSec=600
