# deploy/sports-arb-mlb-fire-collect.timer
[Unit]
Description=Nightly MLB fire-sample collection (08:30 UTC = 4:30 ET)

[Timer]
OnCalendar=*-*-* 08:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
