# deploy/sports-arb-report.service
[Unit]
Description=Sports Arb daily report and Telegram digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/sports-arb
EnvironmentFile=/etc/sports-arb.env
Environment=SPORTS_ARB_STATE_DIR=/var/lib/sports-arb
Environment=SPORTS_ARB_DATA_DIR=/var/lib/sports-arb/data
Environment=SPORTS_ARB_RUNTIME_DIR=/var/lib/sports-arb/runtime
# Root so we can pause the daemon to avoid OOM; npm runs as sports-arb inside the script.
# LLM runs at 04:00 via sports-arb-llm-learn.timer — this job sends report+PnL+Telegram only.
ExecStart=/bin/bash /opt/sports-arb/scripts/nightly-digest.sh telegram
User=root
Group=root
TimeoutStartSec=900
