# deploy/sports-arb-mlb-fire-collect.service
[Unit]
Description=Collect yesterday's MLB wouldFire samples into the backtest repo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=sports-arb
WorkingDirectory=/opt/sports-arb-recorder
EnvironmentFile=/etc/sports-arb-recorder.env
ExecStart=/usr/bin/python3 scripts/collect-mlb-fire-samples.py
