# deploy/sports-arb-backup.timer
[Unit]
Description=Run Sports Arb state backup daily

[Timer]
OnCalendar=*-*-* 00:10:00
Persistent=true
Unit=sports-arb-backup.service

[Install]
WantedBy=timers.target
