# deploy/state-feed-shadow.service
[Unit]
Description=State-feed shadow logger (StatsAPI/FotMob + Polymarket ladders for Strat 1/2)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/sports-arb
EnvironmentFile=/etc/sports-arb.env
Environment=SPORTS_ARB_STATE_DIR=/var/lib/sports-arb
Environment=SPORTS_ARB_DATA_DIR=/var/lib/sports-arb/data
Environment=SPORTS_ARB_RUNTIME_DIR=/var/lib/sports-arb/runtime
Environment=NODE_NO_WARNINGS=1
Environment=STATE_FEED_LIVE_POLL_MS=2000
Environment=STATE_FEED_FOTMOB_POLL_MS=3000
Environment=STATE_FEED_IDLE_POLL_MS=30000
Environment=STATE_FEED_DISCOVERY_MS=60000
ExecStart=/opt/sports-arb/node_modules/.bin/tsx /opt/sports-arb/scripts/state-feed-shadow-logger.ts
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=20
User=sports-arb
Group=sports-arb
StateDirectory=sports-arb
RuntimeDirectory=sports-arb
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
